# elements/rpi-znc/extra-data.d/files/znc.service
#
# ZNC Service configuration generated by rpi-znc element
# Source: https://wiki.znc.in/Running_ZNC_as_a_system_daemon
#
[Unit]
Description=ZNC, an IRC proxy/bouncer
After=network-online.target

[Service]
ExecStart=/usr/bin/znc -f --datadir=/var/lib/znc
User=znc

[Install]
WantedBy=multi-user.target
